# service/dbus-consumption-calculator/run
#!/bin/sh
exec 2>&1
exec /data/dbus-huaweisun2000-pvinverter/dbus-consumption-calculator.py
